# conf/submission/ctf-submission@.service
[Unit]
Description=CTF Flag-Submission Service
Wants=postgresql.service
After=postgresql.service

[Service]
Type=notify
EnvironmentFile=/etc/ctf-gameserver/submission.env
ExecStart=/usr/bin/ctf-submission --port %i
User=nobody
Group=nogroup
RestartSec=10
Restart=on-failure

[Install]
WantedBy=multi-user.target
