# conf/checker/ctf-checkermaster@.service
[Unit]
Description=CTF Checker Runner
After=network.target postgresql.service

[Service]
Type=notify
ExecStart=/usr/bin/ctf-checkermaster @/etc/ctf-gameserver/checkermaster.conf @/etc/ctf-gameserver/checker/%i.conf
ExecStop=/bin/kill -s usr1 $MAINPID
WatchdogSec=20
Restart=on-failure
TimeoutStopSec=65
Environment=PYTHONPATH=/etc/ctf-gameserver/checker/modules/
User=nobody
Group=nogroup
SyslogIdentifier=ctf-checkermaster@%i

[Install]
WantedBy=multi-user.target
