# submission/ctf-submission@.service
[Unit]
Description=CTF Flag-Submission Service

[Service]
Type=notify
ExecStart=/usr/bin/ctf-submission -p %i
User=nobody
Group=nogroup
