# submission/ctf-submission@.service
[Unit]
Description=CTF Flag-Submission Service
Wants=postgresql.service
After=postgresql.service

[Service]
Type=notify
ExecStart=%BINDIR%/ctf-submission @/etc/ctf-gameserver/submission.conf --port %i
User=nobody
Group=nogroup

[Install]
WantedBy=multi-user.target
