# submission/submission.service
[Unit]
Description=CTF Flag-Submission Service

[Service]
Type=notify
ExecStart=flagsubmission
User=nobody
Group=nogroup
