# submission/submission@.service
[Unit]
Description=CTF Flag-Submission Service

[Service]
Type=notify
ExecStart=flagsubmission -p %i
User=nobody
Group=nogroup
