# conf/controller/ctf-scoring.service
[Unit]
Description=FAUST-CTF: update scoring after a game tick
Wants=postgresql.service
After=postgresql.service

[Service]
Type=oneshot
ExecStart=/usr/bin/ctf-scoring
EnvironmentFile=/etc/ctf-gameserver/scoring.env
User=nobody
Group=nogroup
