# conf/controller/ctf-flagid.timer
[Unit]
Description=CTF Gameserver Flagid Exporter

[Timer]
OnCalendar=minutely
AccuracySec=5s

[Install]
WantedBy=timers.target
