# conf/controller/ctf-controller.service
[Unit]
Description=FAUST-CTF: run game tick
Wants=postgresql.service
After=postgresql.service
Before=ctf-scoring.service
Wants=ctf-scoring.service

[Service]
Type=oneshot
ExecStart=/usr/bin/ctf-controller @/etc/ctf-gameserver/controller.conf
User=nobody
Group=nogroup
