# conf/controller/ctf-flagid.service
[Unit]
Description=CTF Gameserver Flagid Exporter
After=postgresql.service

[Service]
Type=oneshot
EnvironmentFile=/etc/ctf-gameserver/flagid.env
ExecStart=/usr/bin/ctf-flagid
User=nobody
Group=nogroup
