# controller/controller.timer
[Unit]
Description=FAUST-CTF: Regularly start game ticks

[Timer]
OnCalendar=*-*-* *:00/03:*