# controller/controller.timer
[Unit]
Description=FAUST-CTF: Regularly start game ticks

[Timer]
OnCalendar=*-*-* *:00/03:00
AccuracySec=5s

[Install]
WantedBy=timers.target
