# controller/ctf-scoring.service
[Unit]
Description=FAUST-CTF: update scoring after a game tick
Wants=postgresql.service
After=postgresql.service

[Service]
Type=oneshot
ExecStart=%BINDIR%/ctf-scoring @/etc/ctf-gameserver/scoring.conf
User=nobody
Group=nogroup
