# controller/ctf-controller.service
[Unit]
Description=FAUST-CTF: run game tick

[Service]
Type=oneshot
ExecStart=/usr/bin/ctf-controller
User=nobody
Group=nogroup
