# checker/checkermaster@.service
[Unit]
Description=CTF Checker Runner

[Service]
Type=notify
ExecStart=checkermaster --config /etc/faustctf/%i
User=nobody
Group=nogroup
