# checker/checkermaster@.service
[Unit]
Description=CTF Checker Runner

[Service]
Type=notify
ExecStart=/usr/bin/checkermaster --config /etc/faustctf/%i.conf --loglevel DEBUG
ExecStop=/bin/kill -USR1 $MAINPID
WatchdogSec=20
Restart=on-failure
TimeoutStopSec=65
Environment=PYTHONPATH=/etc/faustctf/checker
User=nobody
Group=nogroup

[Install]
WantedBy=multi-user.target
